# templates/crontab
SHELL=/bin/bash
0 12 * * * /root/backup-db.sh >> /var/log/cron.log 2>&1
