# templates/crontab
SHELL=/bin/bash
0 0,12 * * * root /root/backup-db.sh >> /var/log/cron.log 2>&1
